# config/launchd/setenv.OLLAMA_HOST.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Sets Ollama environment for all GUI apps (including Ollama.app) at login.
  Only needed on the desktop (the model server).

  Installed by install.sh to ~/Library/LaunchAgents/
  Activate without reboot:  launchctl load ~/Library/LaunchAgents/setenv.OLLAMA_HOST.plist
-->
<dict>
    <key>Label</key>
    <string>setenv.OLLAMA_HOST</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>launchctl setenv OLLAMA_HOST 127.0.0.1:11434; launchctl setenv OLLAMA_MAX_LOADED_MODELS 6</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
